# deploy/factory-scheduler.service
[Unit]
Description=Content Factory — autonomous tick (next approved idea -> run)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Environment=HOME=/home/ubuntu
Environment=CONTENT_AGENTS_DATA=/home/ubuntu/factory-data
Environment=ENGINE_DIR=/home/ubuntu/content-agents
# ключ/base_url/карту тиров наследуем от systemd (без brace-expansion, в отличие от `. .env`)
EnvironmentFile=-/home/ubuntu/factory/.env
ExecStart=/home/ubuntu/factory/factory.sh
